# addon-shelly-emu/rootfs/etc/services.d/shelly-emu/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Shelly Pro 3EM Emulator startet"
exec python3 /app/main.py
